# root/etc/services.d/autossh/run
#!/usr/bin/with-contenv bash

umask 022

#exec \
       s6-setuidgid abc /usr/bin/autossh -M 0 -f -F /config/.ssh/config docker-tunnel-host
